# content/post/linux/2.MakingaService.md~
+++
date = "2015-11-01"
title = "Crear un servicio en Linux"
tags = ["LINUX"]
images = []
menu = "Tags"
banner = "banners/bash.jpg"
categories = ["LINUX"]
+++

Los servicios de linux son útiles a la hora de gestionar daemons del sistema y de hacer que se ejecuten al arrancar el sistema.
<!--more-->

Los servicios son recurson alojados en la carpeta `/etc/init.d/`.

	$vim /etc/init.d/miServicio

Dentro de este archivo escribiremos el siguiente contenido en bash:

	#! /bin/sh

	 ### BEGIN INIT INFO
	# Provides: miServicio
	# Required-Start: $network $local_fs $remote_fs
	# Required-Stop: $network $local_fs $remote_fs
	# Default-Start: 2 3 4 5
	# Default-Stop: 0 1 6
	# Short-Description: miServicio script.
	# Description: Ejemplo de como crear un servio en Linux.
	### END INIT INFO

	NAME=miServicio
	DESC="Demonio de miServicio"
	PIDFILE="/var/run/${NAME}.pid"
	LOGFILE="/var/log/${NAME}.log"

	#ruta del script que invoca el servicio
	DAEMON="/ruta/a/script.sh"

	. /lib/lsb/init-funtions

	do_start(){
		log_daemon_msg "Arrancando demonio $NAME"
		start-stop-daemon --start --background --pidfile $PIDFILE --make-pidfile --startas $DAEMON
		log_end_msg $?
	}

	do_stop () {
    		log_daemon_msg "Matando el demonio $NAME"
	    	start-stop-daemon --stop --pidfile $PIDFILE --retry 10
    		log_end_msg $?
	}

	case "$1" in
    	start)
        	do_start
        	;;
    	stop)
        	do_stop
        	;;
	restart|force-reload)
        	do_stop
		do_start
        	;;
	status)
		status_of_proc "$NAME" "$DAEMON" && exit 0 || exit $?
		;;
    	*)
        	N=/etc/init.d/$NAME
        	echo "Usage: $N {start|stop|restart|force-reload}" >&2
        	exit 1
        	;;
	esac

	exit 0

Por ultimo habrá que darle permisos de ejecución:

	$sudo chmod +x /etc/init.d/miServicio

Si todo ha ido bien ahora podremos iniciar y parar el servicio con tan solo hacer:

	$sudo service miServicio start
	$sudo service miServicio stop
	$sudo service miServicio restart
	$sudo service miServicio status

Con todo esto tenemos el servicio, pero para que se autoejecute al iniciar el sistema tenemos que indicarselo:

	$sudo update-rc.d miServicio defaults
